# config/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noissue --nohostname --autologin ${BCLD_USER} %I ${BCLD_TERM}
Type=idle
